# openstack-auto-scaling.cron
*/10 * * * * stack /bin/bash /usr/bin/openstack-auto-scaling.sh
